# uncensored_creative01/supervisor.conf
[program:ollama]
command=/usr/local/bin/ollama serve
directory=/workspace
user=user
autostart=false
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
environment=HOME="/home/user",OLLAMA_HOST="127.0.0.1:11434",OLLAMA_MODELS="/workspace/models/ollama",OLLAMA_NUM_PARALLEL="1",OLLAMA_MAX_LOADED_MODELS="1",OLLAMA_FLASH_ATTENTION="0"
stdout_logfile=/tmp/xforce-ai/services/ollama/stdout.log
stderr_logfile=/tmp/xforce-ai/services/ollama/stderr.log

[program:ollama-model-init]
command=/opt/creative/bin/init-ollama-model.sh
directory=/workspace
user=user
autostart=false
autorestart=false
startsecs=0
startretries=3
environment=HOME="/home/user",OLLAMA_HOST="127.0.0.1:11434",OLLAMA_MODELS="/workspace/models/ollama"
stdout_logfile=/tmp/xforce-ai/services/ollama-model-init/stdout.log
stderr_logfile=/tmp/xforce-ai/services/ollama-model-init/stderr.log

[program:comfyui]
command=/venv/comfyui/bin/python /opt/creative/ComfyUI/main.py --listen 127.0.0.1 --port 8188 --lowvram --disable-auto-launch --enable-cors-header * --output-directory /workspace/outputs --extra-model-paths-config /opt/creative/comfyui-extra-model-paths.yaml
directory=/opt/creative/ComfyUI
user=user
autostart=false
autorestart=true
startsecs=8
stopasgroup=true
killasgroup=true
environment=HOME="/home/user",PYTHONUNBUFFERED="1",MALLOC_TRIM_THRESHOLD_="65536",LD_PRELOAD="/usr/lib/x86_64-linux-gnu/libtcmalloc.so.4"
stdout_logfile=/tmp/xforce-ai/services/comfyui/stdout.log
stderr_logfile=/tmp/xforce-ai/services/comfyui/stderr.log

[program:creative-portal]
command=/opt/creative/bin/xforce-creative-portal serve --host 127.0.0.1 --port 8090
directory=/workspace
user=user
autostart=true
autorestart=true
startsecs=2
environment=HOME="/home/user",PYTHONUNBUFFERED="1",CREATIVE_PORTAL_STATIC="/opt/creative/portal/static",CREATIVE_FILE_ROOTS="models:/workspace/models,outputs:/workspace/outputs,creative:/workspace/creative"
stdout_logfile=/tmp/xforce-ai/services/creative-portal/stdout.log
stderr_logfile=/tmp/xforce-ai/services/creative-portal/stderr.log
